# symlinks/config/systemd/user/volume.service
[Unit]
Description=Set initial volume and unmute
After=pipewire-pulse.service pulseaudio.service
ConditionFileIsExecutable=%h/.config/volume/init-volume.sh

[Service]
Type=oneshot
RemainAfterExit=yes
# Give PulseAudio/PipeWire time to restore its state, then override it
ExecStartPre=/usr/bin/sleep 3
ExecStart=%h/.config/volume/init-volume.sh

[Install]
WantedBy=default.target
