# symlinks/config/systemd/user/wayland-session.target
[Unit]
Description=Wayland running

BindsTo=graphical-session.target
After=graphical-session.target
