# symlinks/config/systemd/user/xsession.target
[Unit]
Description=xorg running

BindsTo=graphical-session.target
