# env/base/symlinks/config/systemd/user/clean-home-tmp.timer
[Unit]
Description=Clean home tmp timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=basic.target
