# env/base/symlinks/config/systemd/user/clean-home-tmp.timer
[Unit]
Description=Clean ~/tmp every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
