# env/base/symlinks/config/systemd/user/clean-home-tmp.timer
[Unit]
Description=Clean %h/tmp every day

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
