# env/base/symlinks/config/systemd/user/clean-home-tmp.service
[Unit]
Description=Clean home tmp service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "find ~/tmp -type f -mtime +7 -execdir rm -- '{}' \\;"
