# symlinks/config/systemd/user/redshift.service.d/xsession.conf
[Unit]
Wants=xsession.target
After=xsession.target

[Service]
RestartSec=5

[Install]
WantedBy=xsession.target
